# deploy-files/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:laravel-octane-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/alegra/artisan octane:start --server=swoole --host="0.0.0.0" --workers=8 --task-workers=4 --port=8000
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/alegra/storage/logs/main-worker.log
stopwaitsecs=3600

[program:kitchen-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/alegra/artisan queue:work database --name=kitchen --queue=kitchen --tries=3 --timeout=1200
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/alegra/storage/logs/kitchen-worker.log
stopwaitsecs=3600